# Library/LaunchAgents/org.homebrew.ssh-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>org.homebrew.ssh-agent</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/sh</string>
			<string>-c</string>
			<string>rm -f $SSH_AUTH_SOCK; SSH_ASKPASS=/usr/local/bin/ssh-askpass DISPLAY=':0' /opt/homebrew/bin/ssh-agent -D -a $SSH_AUTH_SOCK</string>
		</array>
		<key>RunAtLoad</key>
		<true />
		<key>StandardOutPath</key>
		<string>/tmp/org.homebrew.ssh-agent.out.log</string>
		<key>StandardErrorPath</key>
		<string>/tmp/org.homebrew.ssh-agent.err.log</string>
	</dict>
</plist>
